# ftpserver.rc
#!/bin/sh

. /etc/rc.subr

name=ftpserver
rcvar=ftpserver_enable

command="/scripts/${name}"

load_rc_config $name
run_rc_command "$1"
